# plasma-mobile-wf.service.d/40-droid-wait.conf
[Unit]
BindsTo=android-service@hwcomposer.service
Requires=android-service@hwcomposer.service
After=android-service@hwcomposer.service
