# android-service@hwcomposer.service.d/20-plasma-mobile-wf.conf
[Unit]
BindsTo=plasma-mobile-wf.service
Before=plasma-mobile-wf.service